# xbaydns/tools/master/crontab
0-59/5 * * * * /home/xbaydns/prog/Log2IPlist.sh
0-59/5 * * * * /home/xbaydns/prog/genidcview.sh
0-59/5 * * * * /home/xbaydns/prog/genslaveconf.sh
