# xbaydns/tools/prober/crontab
*/5 * * * * /home/xdagent/prog/UpdateAgent.sh
*/5 * * * * /home/xdagent/prog/UpdateIPList.sh
*/5 * * * * /home/xdagent/prog/ProbeIPlist.sh
*/5 * * * * /home/xdagent/prog/UploadProbeResults.sh
