# xbaydns/tools/slave/crontab
*/5 * * * * /home/xdslave/prog/UpdateNamed.sh
