# xbaydns/tools/slave/crontab
0-59/5 * * * * /home/xdslave/prog/UploadLog.sh
0-59/5 * * * * /home/xdslave/prog/UpdateNamed.sh
